# ops/runtime/launchd/com.ronny.mcp-runtime-anti-drift-cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.mcp-runtime-anti-drift-cycle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ronnyworks/code/agentic-spine/ops/runtime/mcp-runtime-anti-drift-cycle.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/mcp-runtime-anti-drift-cycle.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/mcp-runtime-anti-drift-cycle.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPINE_ROOT</key>
        <string>/Users/ronnyworks/code/agentic-spine</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
</dict>
</plist>
