# ops/runtime/launchd/com.ronny.docker-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronny.docker-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/ssh</string>
    <string>-N</string>
    <string>-i</string>
    <string>/Users/ronnyworks/.ssh/id_ed25519</string>
    <string>-o</string>
    <string>ServerAliveInterval=60</string>
    <string>-o</string>
    <string>ServerAliveCountMax=3</string>
    <string>-o</string>
    <string>ConnectTimeout=10</string>
    <string>-o</string>
    <string>ExitOnForwardFailure=yes</string>
    <string>-o</string>
    <string>StrictHostKeyChecking=no</string>
    <string>-L</string>
    <string>2375:127.0.0.1:2375</string>
    <string>docker-host@100.92.156.118</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SPINE_ROOT</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
    <key>PATH</key>
    <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ronnyworks/code/agentic-spine</string>

  <key>StandardOutPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/docker-tunnel.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/docker-tunnel.err</string>
</dict>
</plist>
