# ops/runtime/launchd/com.ronny.simplefin-daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.simplefin-daily-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>"$HOME/code/workbench/infra/compose/finance/scripts/simplefin-daily-sync.sh"</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/simplefin-daily-sync.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/simplefin-daily-sync.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPINE_ROOT</key>
        <string>/Users/ronnyworks/code/agentic-spine</string>
        <key>PATH</key>
        <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/workbench</string>
</dict>
</plist>
