# ops/runtime/launchd/com.ronny.proposals-retention-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronny.proposals-retention-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ronnyworks/code/agentic-spine/ops/runtime/proposals-retention-daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>40</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SPINE_ROOT</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
    <key>SPINE_OPERATOR_TZ</key>
    <string>America/New_York</string>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>PATH</key>
    <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ronnyworks/code/agentic-spine</string>

  <key>StandardOutPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/proposals-retention-daily.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/proposals-retention-daily.err</string>
</dict>
</plist>
