# ops/runtime/launchd/com.ronny.mailroom-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronny.mailroom-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>set -euo pipefail; SPINE_ROOT="${SPINE_ROOT:-$HOME/code/agentic-spine}"; export SPINE_REPO="$SPINE_ROOT"; TOKEN_FILE="$SPINE_ROOT/mailroom/state/mailroom-bridge.token"; MONITOR_FILE="$SPINE_ROOT/mailroom/state/mailroom-bridge-monitor.token"; [[ -f "$TOKEN_FILE" ]] &amp;&amp; export MAILROOM_BRIDGE_TOKEN="$(tr -d '\r\n' &lt; "$TOKEN_FILE")"; [[ -f "$MONITOR_FILE" ]] &amp;&amp; export MAILROOM_BRIDGE_MONITOR_TOKEN="$(tr -d '\r\n' &lt; "$MONITOR_FILE")"; exec "$SPINE_ROOT/ops/plugins/mailroom-bridge/bin/mailroom-bridge-serve"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SPINE_ROOT</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
    <key>SPINE_OPERATOR_TZ</key>
    <string>America/New_York</string>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>PATH</key>
    <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ronnyworks/code/agentic-spine</string>

  <key>StandardOutPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/mailroom-bridge.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/mailroom-bridge.err</string>
</dict>
</plist>
