# ops/runtime/launchd/com.ronny.media-capacity-snapshot-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.media-capacity-snapshot-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ronnyworks/code/agentic-spine/ops/runtime/media-capacity-snapshot-daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPINE_ROOT</key>
        <string>/Users/ronnyworks/code/agentic-spine</string>
        <key>PATH</key>
        <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SPINE_OPERATOR_TZ</key>
        <string>America/New_York</string>
        <key>TZ</key>
        <string>America/New_York</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/media-capacity-snapshot-daily.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/media-capacity-snapshot-daily.err</string>
</dict>
</plist>
