# ops/runtime/launchd/com.ronny.launchd-health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.launchd-health-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>set -euo pipefail; labels=(com.ronny.freshness-critical-daily com.ronny.verify-run-daily com.ronny.simplefin-daily-sync com.ronny.finance-backup-weekly com.ronny.media-arr-key-roll-drill-weekly); missing=0; for label in "${labels[@]}"; do if launchctl list "$label" >/dev/null 2>/dev/null; then echo "[$(date -u +%Y-%m-%dT%H:%M:%SZ)] OK $label"; else echo "[$(date -u +%Y-%m-%dT%H:%M:%SZ)] MISSING $label"; missing=1; fi; done; exit $missing</string>
    </array>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/launchd-health-out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/launchd-health-err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
</dict>
</plist>
