# ops/runtime/launchd/com.ronny.verify-run-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.verify-run-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/code/agentic-spine/bin/ops cap run verify.run fast</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/verify-run-out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/verify-run-err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPINE_ROOT</key>
        <string>/Users/ronnyworks/code/agentic-spine</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
</dict>
</plist>
