# ops/runtime/launchd/com.ronny.streamdeck.ha.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronny.streamdeck.ha</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>set -euo pipefail; if [[ -f "$HOME/.config/infisical/credentials" ]]; then source "$HOME/.config/infisical/credentials"; fi; HA_API_TOKEN="$($HOME/code/workbench/scripts/agents/infisical-agent.sh get home-assistant prod HA_API_TOKEN)"; export HA_API_TOKEN; exec "$HOME/code/workbench/runtime/streamdeck/.venv/bin/python" "$HOME/code/workbench/runtime/streamdeck/streamdeck_ha_controller.py"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SPINE_ROOT</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
    <key>SPINE_OPERATOR_TZ</key>
    <string>America/New_York</string>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>PATH</key>
    <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ronnyworks/code/workbench/runtime/streamdeck</string>

  <key>StandardOutPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/streamdeck-ha.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/streamdeck-ha.err</string>
</dict>
</plist>
