# ops/runtime/launchd/com.ronny.alerting-probe-cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ronny.alerting-probe-cycle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>set -euo pipefail; SPINE_ROOT="${SPINE_ROOT:-$HOME/code/agentic-spine}"; INF_AGENT="$SPINE_ROOT/ops/tools/infisical-agent.sh"; if [[ -x "$INF_AGENT" ]]; then if [[ -z "${ALERTING_HA_URL:-}" ]]; then ALERTING_HA_URL="$("$INF_AGENT" get-cached infrastructure prod ALERTING_HA_URL 2&gt;/dev/null || true)"; export ALERTING_HA_URL; fi; if [[ -z "${ALERTING_HA_TOKEN:-}" ]]; then ALERTING_HA_TOKEN="$("$INF_AGENT" get-cached infrastructure prod ALERTING_HA_TOKEN 2&gt;/dev/null || true)"; export ALERTING_HA_TOKEN; fi; fi; exec "$SPINE_ROOT/ops/runtime/alerting-probe-cycle.sh"</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/alerting-probe-cycle.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/ronnyworks/code/agentic-spine/mailroom/logs/alerting-probe-cycle.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SPINE_ROOT</key>
        <string>/Users/ronnyworks/code/agentic-spine</string>
        <key>SPINE_OPERATOR_TZ</key>
        <string>America/New_York</string>
        <key>TZ</key>
        <string>America/New_York</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
</dict>
</plist>
