# ops/runtime/launchd/com.ronny.mailroom-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ronny.mailroom-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ronnyworks/code/agentic-spine/ops/plugins/mailroom-bridge/bin/mailroom-bridge-serve</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SPINE_REPO</key>
    <string>/Users/ronnyworks/code/agentic-spine</string>
    <key>SPINE_INBOX</key>
    <string>/Users/ronnyworks/code/.runtime/spine-mailroom/inbox</string>
    <key>SPINE_OUTBOX</key>
    <string>/Users/ronnyworks/code/.runtime/spine-mailroom/outbox</string>
    <key>SPINE_STATE</key>
    <string>/Users/ronnyworks/code/.runtime/spine-mailroom/state</string>
    <key>SPINE_LOGS</key>
    <string>/Users/ronnyworks/code/.runtime/spine-mailroom/logs</string>
    <key>MAILROOM_BRIDGE_TOKEN_FILE</key>
    <string>/Users/ronnyworks/code/.runtime/spine-mailroom/state/mailroom-bridge.token</string>
    <key>SPINE_OPERATOR_TZ</key>
    <string>America/New_York</string>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>PATH</key>
    <string>/Users/ronnyworks/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/ronnyworks/code/agentic-spine</string>

  <key>StandardOutPath</key>
  <string>/Users/ronnyworks/code/.runtime/spine-mailroom/logs/mailroom-bridge.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/ronnyworks/code/.runtime/spine-mailroom/logs/mailroom-bridge.err</string>
</dict>
</plist>
